# oauth2_proxy/supervisor.conf
[program:oauth2_proxy]
command=oauth2_proxy --upstream="http://127.0.0.1:80" --http-address="0.0.0.0:4180"
