# atlantis/supervisor.conf
[program:atlantis]
command=atlantis server --repo-config="/etc/atlantis/repos.yaml" --port=4141
